# menu.service
[Unit]
Description=Menu Flask Application
Documentation=https://github.com/your-username/menu
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/menu
Environment=PATH=/var/www/menu/venv/bin
Environment=PYTHONPATH=/var/www/menu
Environment=FLASK_ENV=production
Environment=FLASK_APP=app.py

# 启动命令
ExecStart=/var/www/menu/venv/bin/python app.py

# 重启策略
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 安全配置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/menu
ReadWritePaths=/var/www/menu/instance
ReadWritePaths=/var/www/menu/static/images

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=menu

# 健康检查
ExecStartPost=/bin/sleep 5
ExecStartPost=/bin/bash -c 'curl -f http://localhost:8081/health || exit 1'

[Install]
WantedBy=multi-user.target
